# docker-cassandra-spark-slave/supervisord.conf
[supervisord]
nodaemon=true

[program:sparkslave]
command=/app/spark-2.2.0-bin-hadoop2.7/sbin/start-slave.sh spark://sparkmaster:7077

[program:cassandra]
command=/docker-entrypoint.sh cassandra -f
