# docker-spark-master/supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/app/spark-2.2.0-bin-hadoop2.7/sbin/start-master.sh 
